# config/.config/systemd/user/check-updates.service
[Unit]
Description=Check for system and AUR updates

[Service]
Type=oneshot
ExecStart=%h/.config/K480X/scripts/check-updates.sh

